# rootfs/etc/services.d/cantaloupe/run
#!/usr/bin/with-contenv sh

exec /usr/bin/java -Dcantaloupe.config=/usr/local/cantaloupe/cantaloupe.properties -Xmx8g -jar /usr/local/cantaloupe/Cantaloupe-3.4.3.war