# rootfs/etc/services.d/cantaloupe/run
#!/usr/bin/with-contenv sh

exec java $JAVA_OPTS -Dcantaloupe.config=/usr/local/cantaloupe/cantaloupe.properties -jar /usr/local/cantaloupe/Cantaloupe-3.4.3.war