# rootfs/etc/services.d/cantaloupe/run
#!/usr/bin/with-contenv sh

exec java $JAVA_OPTS -Dcantaloupe.config=$CANTALOUPE_HOME/cantaloupe.properties -jar $CANTALOUPE_HOME/Cantaloupe-3.4.3.war